# DiscordBot/discord.service
[Unit]
Description=discord bot

[Service]
ExecStart=/usr/bin/python3 /home/ec2-user/app/discordbot.py

[Install]
WantedBy=multi-user.target
